# _devops/certificate_renewal/crontab
5 4 * * *    root    /scripts/renew-certificates >> /var/log/cron.log 2>&1
