# _devops/certificate_renewal/crontab
# min    hour    day    month    weekday    command
5        4       *      *        *          /scripts/renew-certificates
